# consul/supervisor.conf
[program:consul]
command=consul agent -data-dir=/consul/data -config-dir=/etc/consul
